# infra/systemd/cloud-reboot.path
[Unit]
Description=Watch for the deniz-cloud reboot sentinel

[Path]
PathExists=/var/lib/deniz-cloud/reboot-requested
Unit=cloud-reboot.service

[Install]
WantedBy=multi-user.target
